# mpd-notifier/image_server.service
[Service]
ExecStart=/usr/bin/python3 /path/to/image_server.py 9555
Restart=always
RestartSec=3

